# openwrt-18.06/target/linux/siflower/sf19a28-fullmask/base-files-SF19A28-AC22/etc/hotplug.d/iface/80-updateserverip
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0
[ "$INTERFACE" = "wwan" ] || [ "$INTERFACE" = "wan" ] || exit 0
aclscript l_update_ip &
echo "update lease network server ip" > /dev/ttyS0
