# openwrt-18.06/target/linux/siflower/sf19a28-fullmask/base-files-SF19A28-AC28/etc/hotplug.d/iface/80-updatewan
#!/bin/sh

update_hnat_wanmask() {
	wan_masklen=`ubus call network.interface.wan status | grep 'ipv4-address' -A 6 | grep mask | awk -F ' ' '{print $2}'| awk -F ',' '{print $1}' | head -n 1 |tr -d "\n"`
	if [ "x$wan_masklen" != "x" -a "$wan_masklen" != "0" ]; then
		if [ "$wan_masklen" -ge 32 ]; then
			wan_masklen=31
		fi
		echo "updatewanmask $wan_masklen" > /sys/kernel/debug/hnat_debug
	fi
}

[ "$ACTION" = ifup ] || exit 0
[ "$INTERFACE" = "wwan" ] || [ "$INTERFACE" = "wan" ] || exit 0
update_hnat_wanmask
