# openwrt-18.06/target/linux/siflower/sf19a28-fullmask/base-files-SF19A28-AC28/etc/hotplug.d/iface/80-updatewan
#!/bin/sh

update_hnat_wanmask() {
	wan_masklen=`ubus call network.interface.wan status | grep 'ipv4-address' -A 6 | grep mask | awk -F ' ' '{print $2}'| awk -F ',' '{print $1}' | head -n 1 |tr -d "\n"`
	if [ "x$wan_masklen" != "x" -a "$wan_masklen" != "0" ]; then
		if [ "$wan_masklen" -ge 32 ]; then
			wan_masklen=31
		fi
		echo "updatewanmask $wan_masklen" > /sys/kernel/debug/hnat_debug
	fi
}

update_wan_ip() {
	wan_ip=`ubus call network.interface.wan status | grep 'ipv4-address' -A 6 | grep "\"address\"\:" | awk -F '"' '{print $4}' | head -n 1 | tr -d "\n"`
	old_ip=`uci -q get network.wan.recent_ip`
	if [ "$wan_ip" != "$old_ip" -a "x$wan_ip" != "x" ]; then
		if [ "x$old_ip" != "x" ]; then
			echo "vldclean" > /sys/kernel/debug/hnat_debug
		fi
		uci set network.wan.recent_ip=$wan_ip
		uci commit network
	fi
}

[ "$ACTION" = ifup ] || exit 0
[ "$INTERFACE" = "wwan" ] || [ "$INTERFACE" = "wan" ] || exit 0
update_hnat_wanmask
update_wan_ip
