# openwrt-18.06/target/linux/siflower/sf19a28-fullmask/base-files-SF19A28-AC28/etc/hotplug.d/iface/。
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

if echo "$INTERFACE" | grep -q "lan"
then
	echo "$INTERFACE is in lan"
else
	exit 0
fi

iptohex() {
	IFS=.
	for str in $1
	do
		printf "%02x" $str
	done
}

echo clearlan > /sys/kernel/debug/hnat_debug

for i in $(seq 1 7)
do
	if [ $i == 1 ]; then
		lan_name=lan
	else
		lan_name=lan$i
	fi
	lanip=`uci -q get network.$lan_name.ipaddr`
	if [ "x$lanip" != "x" ]; then
		lanhex=`iptohex $lanip`
		netmask=`uci -q get network.$lan_name.netmask`
		if [ "x$netmask" == "x" ]; then
			netmask="255.255.255.0"
		fi
		maskhax=`iptohex $netmask`
		num=$(($i-1))
		echo "setlan $num 0x$lanhex 0x$maskhax" > /sys/kernel/debug/hnat_debug
	fi
done

echo updatelan > /sys/kernel/debug/hnat_debug
